# vietanhbui2000/luci-app-ipkinstaller/etc/uci-defaults/luci-app-ipkinstaller
#!/bin/sh

sed -i 's/cbi.submit\"] = true/cbi.submit\"] = \"1\"/g' /usr/lib/lua/luci/dispatcher.lua

rm -f /tmp/luci-app-filebrowser
exit 0
