# vietanhbui2000/luci-theme-advancedtomatomaterial/root/etc/uci-defaults/luci-theme-advancedtomatomaterial
#!/bin/sh

uci batch <<-EOF
	set luci.themes.advancedtomatomaterial=/luci-static/advancedtomatomaterial
	commit luci
EOF

exit 0
